# crontab
# Cron файл для запуска программы каждые 10 минут
*/10 * * * * /app/run.sh >> /app/cron.log 2>&1
